# systemd/user/LastfmMonthly.timer
[Unit]
Description=Monthly Lastfm Collage Timer

[Timer]
OnCalendar=*-*~01 21:00:00
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
